# src/plonex/supervisor/templates/program.conf.j2
[program:{{ program }}]
command = {{ command }}
process_name = {{ process_name }}
directory = {{ directory }}
priority = {{ priority }}
redirect_stderr = false
stopasgroup = true
